# test/fixtures/strong-pm-1.4.conf
# strong-pm
# Generated by strong-service-upstart on 2014-08-22T00:00:00.000Z

description "Node application server"
author "Current User"

# Start on boot/reboot
start on runlevel [2345]

# Stop on shutdown/reboot
stop on runlevel [016]

# To start/stop strong-pm service manually, run:
#
# /sbin/initctl start strong-pm
# /sbin/initctl stop strong-pm

respawn

# Log all service output to /var/log/strong-pm.log
# Log rotation is handled by Upstart
console log

# User & Group to run service as:
chdir /
setuid nobody
setgid nogroup

# Arbitrarily high limit on number of open file descriptors:
limit nofile 50000 50000

# Enable core dumps:
limit core unlimited unlimited

# Environment variables:
# no environment variables set

# Run our actual application
exec /custom/path/to/node /custom/path/to/my/script.js
